# openwrt/files/wifi-portal.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99

BIN=/usr/bin/wifi-portal

start_service() {
	procd_open_instance
	procd_set_param command $BIN
	procd_close_instance
}

